# NodeExporter/Linux/node_exporter.service
[Unit]
Description=Prometheus Node Exporter
Documentation=https://github.com/prometheus/node_exporter
After=network-online.target
Wants=network-online.target
# StartLimitInterval must be > StartLimitBurst * RestartSec
StartLimitInterval=60
StartLimitBurst=5

[Service]
User=node_exporter
Group=node_exporter
Type=simple
SyslogIdentifier=node_exporter

ExecStart=/usr/bin/node_exporter
# To listen on a custom port:
# ExecStart=/usr/bin/node_exporter --web.listen-address=:9200

Restart=on-failure
RestartSec=5
TimeoutStartSec=30
TimeoutStopSec=15

# Security hardening
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ProtectControlGroups=yes
ProtectKernelModules=yes
ProtectKernelTunables=yes
RestrictAddressFamilies=AF_INET AF_INET6
RestrictNamespaces=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
MemoryDenyWriteExecute=yes
LockPersonality=yes
SystemCallArchitectures=native
# node_exporter runs on port 9100 and needs no special capabilities
CapabilityBoundingSet=

# Process limits
LimitNOFILE=65535
LimitNPROC=4096

[Install]
WantedBy=multi-user.target
