# NodeExporter/Linux/node_exporter.service
[Unit]
Description=Prometheus Node Exporter
Documentation=https://github.com/prometheus/node_exporter
After=network-online.target
Wants=network-online.target
# StartLimitInterval must be > StartLimitBurst * RestartSec
StartLimitInterval=60
StartLimitBurst=5

[Service]
User=node_exporter
Group=node_exporter
Type=simple
SyslogIdentifier=node_exporter

ExecStart=/usr/bin/node_exporter
# To listen on a custom port:
# ExecStart=/usr/bin/node_exporter --web.listen-address=:9200

Restart=on-failure
RestartSec=5
TimeoutStartSec=30
TimeoutStopSec=15

# Security hardening
# Keep it light: node_exporter needs read access to /proc, /sys, and
# /sys/fs/cgroup to export host-level metrics (network interfaces, disks,
# kernel tunables, cgroups, etc.).
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=full
ProtectHome=yes
ProtectKernelModules=yes
RestrictAddressFamilies=AF_INET AF_INET6 AF_NETLINK
RestrictRealtime=yes
RestrictSUIDSGID=yes
LockPersonality=yes
SystemCallArchitectures=native
CapabilityBoundingSet=CAP_DAC_READ_SEARCH CAP_NET_ADMIN CAP_SYS_PTRACE
AmbientCapabilities=CAP_DAC_READ_SEARCH CAP_NET_ADMIN CAP_SYS_PTRACE

# Process limits
LimitNOFILE=65535
LimitNPROC=4096

[Install]
WantedBy=multi-user.target
