# Scaphandre/Linux/scaphandre.service
[Unit]
Description=Scaphandre Prometheus power metrics exporter
Documentation=https://hubblo-org.github.io/scaphandre-documentation/
After=network-online.target
Wants=network-online.target
# StartLimitInterval must be > StartLimitBurst * RestartSec
StartLimitInterval=60
StartLimitBurst=5

[Service]
# Scaphandre requires root to read RAPL power data from /sys/class/powercap/
User=root
Group=root
Type=simple
SyslogIdentifier=scaphandre

ExecStart=/usr/bin/scaphandre prometheus \
    --qemu \
    --containers \
    --port 1920 \
    --address 0.0.0.0

Restart=on-failure
RestartSec=5
TimeoutStartSec=30
TimeoutStopSec=15

# Security hardening (limited by root requirement, but still useful)
NoNewPrivileges=yes
PrivateTmp=yes
ProtectHome=yes
ProtectKernelModules=yes
RestrictAddressFamilies=AF_INET AF_INET6
RestrictRealtime=yes
LockPersonality=yes
SystemCallArchitectures=native

# Process limits
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
